# supervisord-m0.conf
[supervisord]
loglevel = debug
nodaemon = false

[program:masternode]
command=npm start
