# supervisord-l0.conf
[supervisord]
loglevel = debug
nodaemon = false

[program:lightnode]
command=npm start
